# smokeping/root/etc/supervisord/conf.d/smokeping.conf
[program:smokeping]
command=/usr/sbin/smokeping
autostart=true
autorestart=true
stderr_logfile=/var/log/smokeping.err.log
stdout_logfile=/var/log/smokeping.out.log
